# redis_gpsd_ubxtool/cgn_services/cgn_ntrip.service
[Unit]
Description= --- NtripClient trigger ---
After=multi-user.target

[Service]
Type=exec
ExecStart=/usr/bin/python3 /home/agrodroid/releases/active_release/redis_gpsd_ubxtool/ntrip.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
